# systemd/.config/systemd/user/zita-j2n-restart-trigger.path
[Unit]
Wants=zita-j2n.service

[Path]
PathModified=/home/johannes/.config/environment.d/10-zita-j2n.conf
